# deploy/ec2/appseed.service
# /etc/systemd/system/appseed.service
# Systemd service for AppSeed docker-compose auto-restart
#
# Install:
#   sudo cp appseed.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now appseed

[Unit]
Description=AppSeed Docker Compose Application
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/appseed

# Pull latest images and start all prod services
ExecStart=/usr/bin/docker compose --profile prod up -d --build --remove-orphans
ExecStop=/usr/bin/docker compose --profile prod down
ExecReload=/usr/bin/docker compose --profile prod up -d --build --remove-orphans

# Restart policy — systemd restarts the compose stack if it exits unexpectedly
Restart=on-failure
RestartSec=10

# Environment file for secrets (optional, docker-compose.yml reads .env by default)
EnvironmentFile=-/opt/appseed/.env

[Install]
WantedBy=multi-user.target
